# imageroot/systemd/user/2fauth.service
#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

#
# This systemd unit starts a 2fauth instance using Podman.
# Most parts of this file come from podman-generate-systemd.
#

[Unit]
Description=Podman 2fauth.service
Requires=2fauth-app.service
Before=2fauth-app.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=-%S/state/environment
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/2fauth.pid %t/2fauth.pod-id
ExecStartPre=/usr/bin/podman pod create --infra-conmon-pidfile %t/2fauth.pid \
    --pod-id-file %t/2fauth.pod-id \
    --name 2fauth \
    --publish 127.0.0.1:${TCP_PORT}:80 \
    --replace
ExecStart=/usr/bin/podman pod start --pod-id-file %t/2fauth.pod-id
ExecStop=/usr/bin/podman pod stop --ignore --pod-id-file %t/2fauth.pod-id -t 10
ExecStopPost=/usr/bin/podman pod rm --ignore -f --pod-id-file %t/2fauth.pod-id
PIDFile=%t/2fauth.pid
Type=forking

[Install]
WantedBy=default.target
